# protected/extensions/crontab/crontabs/my_crontab
* * * * * php C:/xampp/htdocs/mts/yiicmd.php test 'datetime'
0 0 * * * php C:/xampp/htdocs/mts/yiicmd.php test dor
